# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
childlogdir=/tmp

; Telegram / channels agent loop
[program:nanobot-gateway]
command=nanobot gateway
autostart=true
autorestart=true
startretries=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
priority=10

; OpenAI-compatible HTTP API (used by OpenSpace dashboard)
[program:nanobot-serve]
command=nanobot serve --host 0.0.0.0 --port 18790
autostart=true
autorestart=true
startretries=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
priority=20
